# deploy/econext-gateway.service
# Reference systemd unit for ecoNEXT Gateway.
# The install.sh script generates the actual unit with correct paths.
# To install manually, copy this file to /etc/systemd/system/ and
# adjust User, WorkingDirectory, ExecStart, and Environment values.

[Unit]
Description=ecoNEXT Gateway
After=network.target
Wants=network.target

[Service]
Type=simple
User=linus
WorkingDirectory=/opt/econext-gateway
ExecStart=/opt/econext-gateway/.venv/bin/uvicorn econext_gateway.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StateDirectory=econext-gateway
Environment=ECONEXT_SERIAL_PORT=/dev/econext
Environment=ECONEXT_LOG_LEVEL=INFO

[Install]
WantedBy=multi-user.target
